# systemdmode/quadlets/reloader.container
# reloader.container
# OpenPERouter FRR Configuration Reloader
#
# Sidecar container that receives FRR configuration updates via Unix socket
# and applies them to the FRR daemon without pod restarts.

[Unit]
Description=OpenPERouter FRR Reloader Container
Documentation=https://github.com/openperouter

[Container]
# Container image - same as FRR container
Image=quay.io/openperouter/router:main

# Pod association - runs in same pod as FRR container
Pod=routerpod.pod

# Container name within the pod
ContainerName=reloader

# Volume mounts
# Named volume for FRR sockets - shared with FRR container
Volume=frr-sockets.volume:/var/run/frr:Z

# Host directory mount for FRR configuration
# :Z = Private SELinux label
# Controller writes config here, reloader reads and applies to FRR
Volume=/etc/perouter/frr:/etc/perouter:Z

# Override container entrypoint to run reloader binary
Entrypoint=/reloader

# Reloader command-line arguments
# --frrconfig    - Path to FRR configuration file (written by controller)
# --loglevel     - Logging verbosity
# --unixsocket   - Unix socket path for receiving reload requests from controller
Exec=--frrconfig=/etc/perouter/frr.conf --loglevel=debug --unixsocket /etc/perouter/frr.socket

# Health check - kill container if FRR daemons stop responding,
# systemd Restart=on-failure handles restart
HealthCmd=wget -qO- http://127.0.0.1:9080/healthz || exit 1
HealthInterval=20s
HealthTimeout=5s
HealthRetries=3
HealthOnFailure=kill

[Service]
# Restart policy - restart container if it fails
Restart=on-failure

# Timeout for stopping the container
TimeoutStopSec=70

[Install]
# Automatically enable this container at boot
WantedBy=default.target
